# userpatches/overlay/lib/systemd/system/iperf3.service
[Unit]
Description=iperf3 server
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/iperf3 -s -i1
